# yapi/backup.conf.cron
#!/bin/bash

* * * * * ~/devOps/yapi/mongo/backup.sh
